# files/etc/cron.d/certbot
# Run certbot scripts
0 0 * * * root /usr/local/bin/updatecerts.sh > /dev/null
